# debian/cinder-volume.postinst.in
#!/bin/sh

set -e

#PKGOS-INCLUDE#

if [ "$1" = "configure" ] || [ "$1" = "reconfigure" ] ; then
	pkgos_var_user_group cinder
	chmod 0750 /var/lib/cinder/volumes
	chown root:root /etc/cinder/rootwrap.d/volume.filters

	if which invoke-rc.d >/dev/null 2>&1; then
		invoke-rc.d tgtd restart > /dev/null 2>&1 || true
	else
		/etc/init.d/tgtd restart > /dev/null 2>&1 || true
	fi
fi

#DEBHELPER#

exit 0
