# debian/cinder-volume.service
[Unit]
Description=OpenStack Cinder Volume
Documentation=http://docs.openstack.org/
After=postgresql.service mysql.service keystone.service rabbitmq-server.service ntp.service network.service

[Service]
User=cinder
Group=cinder
ExecStart=/etc/init.d/cinder-volume systemd-start
ExecStop=/etc/init.d/cinder-volume systemd-stop
PIDFile=/var/run/cinder/cinder-volume.pid
Restart=on-failure
Type=forking

[Install]
WantedBy=multi-user.target

