# debian/cinder-scheduler.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d cinder-scheduler remove
fi

#DEBHELPER#
